# build/server-dlmgr-deluge/files/deluge-web.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
  Assumes user=dlmgr group=dlmgr
  Assumes '{{PREFIX}}/dlmgr/deluge' is installation directory

  To enable/disable
  svcadm enable deluged
  svcadm disable deluged

  To check logs
  tail /var/svc/log/network-deluged\:default.log
-->
<service_bundle type='manifest' name='deluge-web'>
  <service name='network/dlmgr/deluge-web' type='service' version='1'>
    <create_default_instance enabled='false' />
    <single_instance />

    <!-- Only start in muti-user mode -->
    <dependency name='multi-user' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/multi-user' />
    </dependency>

    <!-- Wait for network interfaces to be initialized. -->
    <dependency name='network' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/network:default' />
    </dependency>

    <!-- Wait for all local filesystems to be mounted. -->
    <dependency name='filesystem-local' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local:default' />
    </dependency>

    <!-- execute as user USERNAME -->
    <method_context>
      <method_credential user='dlmgr' group='dlmgr' />
      <method_environment>
        <envvar name="PATH" value="/opt/obd/bin:/opt/omni/bin:/usr/gnu/bin:/usr/bin:/bin" />
      </method_environment>
    </method_context>

    <exec_method type='method' name='start' exec='/usr/bin/deluge-web -f -c {{PREFIX}}/dlmgr/.config/deluge -l {{PREFIX}}/dlmgr/.config/deluge/deluge-web.log -L warning' timeout_seconds='60' />
    <exec_method type='method' name='stop' exec=':kill' timeout_seconds='2' />

    <property_group name="startd" type="framework">
      <propval name="duration" type="astring" value="contract" />
      <propval name="ignore_error" type="astring" value="core,signal" />
    </property_group>

    <template>
      <common_name>
        <loctext xml:lang='C'>deluge-web</loctext>
      </common_name>
      <documentation>
        <doc_link name='deluge'  uri='http://dev.deluge-torrent.org/' />
        <manpage title='deluge-web' section='1M' manpath='/usr/share/man' />
      </documentation>
    </template>
  </service>
</service_bundle>
