# build/server-krb5-heimdal/files/smf/heimdal/kdc.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='network/security/heimdal/kdc' type='service' version='0'>
    <create_default_instance enabled='false'/>
    <single_instance/>
    <dependency name='dns' grouping='require_all' restart_on='error' type='service'>
      <service_fmri value='svc:/network/dns/client'/>
    </dependency>
    <exec_method name='start' type='method' exec='{{PREFIX}}/sbin/kdc' timeout_seconds='60'>
      <method_context>
        <method_credential user='root' group='root' privileges='basic,!file_link_any,!proc_info,!proc_session,net_privaddr,proc_audit'/>
      </method_context>
    </exec_method>
    <exec_method name='stop' type='method' exec=':kill' timeout_seconds='60'/>
    <stability value='Unstable'/>
    <template>
      <common_name>
        <loctext xml:lang='C'>Kerberos key distribution center</loctext>
      </common_name>
      <documentation>
        <manpage title='kdc' section='1M' manpath='{{PREFIX}}/share/man'/>
      </documentation>
    </template>
  </service>
</service_bundle>
