# build/download-manager-sabnzbd/files/smf.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
  Assumes user=dlmgr group=dlmgr
  Assumes '/opt/obd/dlmgr/sabnzbd' is installation directory

  To enable/disable
    svcadm enable sabnzbd
    svcadm disable sabnzbd

  To check logs
    tail /var/svc/log/network-sabnzbd\:default.log
-->

<service_bundle type='manifest' name='sabnzbd'>

  <service
    name='network/sabnzbd'
    type='service'
    version='1'>

    <create_default_instance enabled='false' />
    <single_instance />

   <!-- Only start in muti-user mode -->
   <dependency name='multi-user'
      grouping='require_all'
      restart_on='none'
      type='service'>
      <service_fmri value='svc:/milestone/multi-user' />
   </dependency>

    <!-- Wait for network interfaces to be initialized. -->
    <dependency name='network'
        grouping='require_all'
        restart_on='none'
        type='service'>
        <service_fmri value='svc:/milestone/network:default'/>
    </dependency>

    <!-- Wait for all local filesystems to be mounted. -->
    <dependency name='filesystem-local'
        grouping='require_all'
        restart_on='none'
        type='service'>
        <service_fmri value='svc:/system/filesystem/local:default'/>
    </dependency>

    <!-- execute as user USERNAME -->
    <method_context>
        <method_credential user='dlmgr' group='dlmgr' />
	<method_environment>
		<envvar name="PATH" value="/opt/omni/bin:/usr/gnu/bin:/usr/bin:/bin" />
	</method_environment>
    </method_context>

    <exec_method
        type='method'
        name='start'
        exec='/usr/bin/python /opt/obd/dlmgr/sabnzbd/SABnzbd.py -d -f /opt/obd/dlmgr/.config/sabnzbd/sabnzbd.ini'
        timeout_seconds='60'>
    </exec_method>

    <exec_method
        type='method'
        name='stop'
        exec=':kill'
        timeout_seconds='2'>
    </exec_method>

    <property_group name="startd" type="framework">
        <propval name="duration" type="astring" value="contract" />
        <propval name="ignore_error" type="astring" value="core,signal" />
    </property_group>

    <template>
       <common_name>
          <loctext xml:lang='C'>SABnzbd</loctext>
       </common_name>
       <documentation>
          <doc_link name='SABnzbd'  uri='www.http://sabnzbd.org/' />
       </documentation>
    </template>

  </service>

</service_bundle>
