# build/mdns-announcer/staging/_rootfs_/lib/svc/manifest/network/mdns-announcer.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
 
<service_bundle type='manifest' name='mdns-announcer'>
 
  <service
    name='network/mdns-announcer'
    type='service'
    version='1'>
 
    <create_default_instance enabled='false' />
    <single_instance />
  
    <!-- Wait for dns mulitcast to be initialized. -->
    <dependency name='network'
        grouping='require_all'
        restart_on='none'
        type='service'>
        <service_fmri value='svc:/network/dns/multicast:default'/>
    </dependency>
 
    <!-- Wait for all local filesystems to be mounted. -->
    <dependency name='filesystem-local'
        grouping='require_all'
        restart_on='none'
        type='service'>
        <service_fmri value='svc:/system/filesystem/local:default'/>
    </dependency>
 
    <!-- execute as user USERNAME -->
    <method_context>
        <method_credential user='mdns' group='mdns' />
    </method_context>

    <exec_method type='method' name='start' exec='/usr/sbin/mdns-announcer' timeout_seconds='60' />
    <exec_method type='method' name='stop'  exec=':kill' timeout_seconds='60' />

    <property_group name="startd" type="framework">
        <propval name="duration" type="astring" value="contract" />
        <propval name="ignore_error" type="astring" value="core,signal" />
    </property_group>

    <template>
       <common_name>
          <loctext xml:lang='C'>mdnsd</loctext>
       </common_name>
    </template>
 
  </service>
 
</service_bundle>
